# 01-Software-Engineering-Best-Practices/01-Setup/02-Linux-and-Bash/check_ssh.service
[Unit]
Description=Checking connection

[Service]
ExecStart=/bin/bash /usr/local/check_ssh.sh

[Install]
WantedBy=multi-user.target
